# services/telegram/deploy/com.founders.telegram.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.founders.telegram.bot</string>
	<key>ProgramArguments</key>
	<array>
		<string>REPLACE_WITH_VAULT_ROOT/services/telegram/deploy/run-bot.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>REPLACE_WITH_VAULT_ROOT/services/telegram</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>REPLACE_WITH_HOME</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>REPLACE_WITH_HOME/Library/Logs/founders-telegram/bot.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>REPLACE_WITH_HOME/Library/Logs/founders-telegram/bot.stderr.log</string>
</dict>
</plist>
